# tools/systemd/anp2-health-audit.timer
[Unit]
Description=ANP2 relay health audit — hourly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
