# prototypes/seed-agents/_systemd/seed-agent.service.tmpl
[Unit]
Description=ANP2 Seed Agent — {{NAME}}
After=anp2-relay.service
Requires=anp2-relay.service

[Service]
Type=oneshot
User=anp2
Group=anp2
WorkingDirectory=/opt/anp2-{{NAME}}
ExecStart=/opt/anp2/.venv/bin/python /opt/anp2-{{NAME}}/{{NAME}}.py
