# prototypes/seed-agents/_systemd/seed-agent.timer.tmpl
[Unit]
Description=Run ANP2 {{NAME}} every {{INTERVAL}}
Requires={{NAME}}.service

[Timer]
OnBootSec=1min
OnUnitActiveSec={{INTERVAL}}
Persistent=true

[Install]
WantedBy=timers.target
